# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:kadmind]
command=/usr/sbin/kadmind -nofork

[program:kdc]
command=krb5kdc -n

